# deploy/systemd/fiber-ent-apollo-pg.socket
[Unit]
Description=Socket for Fiber Ent Apollo PG App

[Socket]
ListenStream=0.0.0.0:8080

[Install]
WantedBy=sockets.target

